# vendor/QITAS
#!/bin/sh
#chkconfig: 2345 80 90
    #description:qitas

if [ -f /usr/local/vendor.ko ]; then
  sudo insmod /usr/local/qitas.ko
fi

exit 0
