# info/QITAS
#!/bin/sh
#chkconfig: 2345 80 90
    #description:rex

if [ -f /usr/local/vendor.ko ]; then
  sudo insmod /usr/local/vendor.ko
fi

sudo echo "finish loading rex vendor info"
exit 0
