# mark/QITAS
#!/bin/sh
#chkconfig: 2345 80 90
    #description:qitas

if [ -f /usr/local/vendor.ko ]; then
  sudo insmod /usr/local/vendor.ko
fi

sudo echo "finish loading qitas mark config"
exit 0
